# redirect-center.service
[Unit]
Description=Redirect Center - DNS CNAME redirect service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/redirect-center
ExecStart=/usr/bin/deno run --allow-net --allow-read --allow-env --unstable-kv src/main.ts
Restart=always
RestartSec=3
Environment=DENO_DIR=/opt/redirect-center/.deno

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redirect-center

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
