# redirect-center.service
[Unit]
Description=Redirect Center - DNS CNAME redirect service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/redirect.center
ExecStart=/usr/bin/deno run --allow-net --allow-read --allow-env --v8-flags=--max-old-space-size=512 src/main.ts
Restart=always
RestartSec=3
Environment=FQDN=redirect.center
Environment=ENTRY_IP=54.84.55.102
Environment=LISTEN_PORT=80
Environment=LOGGER_LEVEL=info

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redirect-center

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
